# setup/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --dbpath /data/db
autorestart=true

[program:nodeServer]
command=/usr/bin/node /usr/src/app/index.js
autorestart=true
